# marzforwarder-renew.timer
[Unit]
Description=Daily SSL renewal timer for MarzForwarder

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
